# mail/var/spool/cron/crontabs/root
0 5 * * 0 postfix reload && dovecot reload
