# filter/var/spool/cron/crontabs/root
0 1 * * * /usr/bin/freshclam
0 2 * * * /usr/bin/sa-update
0 * * * * /usr/bin/sa-learn --dbpath /var/lib/amavis/.spamassassin --spam /home/virtual/*/*/Maildir/.Junk/cur
0 * * * * /usr/bin/sa-learn --dbpath /var/lib/amavis/.spamassassin --ham /home/virtual/*/*/Maildir/cur
