# filter/var/spool/cron/crontabs/root
0 * * * * /usr/bin/freshclam --quiet
0 * * * * /usr/bin/sa-update
